# android9/overlay/lib/systemd/system/mount-android.service
[Unit]
Description=Remount / rw
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/sbin/mount-android.sh

[Install]
WantedBy=sysinit.target
