# android9/overlay/etc/systemd/system/phosh.service.d/20-wlroots-hwc-env.conf
[Service]
Environment=WLR_BACKENDS='hwcomposer,libinput'
